# exe/controller/thermostat_sensor.service
[Unit]
Description=ThermostatService
After=thermostat_controller.service

[Service]
ExecStart=python3 mainSensor.py
WorkingDirectory=/home/pi/Documents/Thermostat/exe/controller
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
